# blackwtredis/crontab
0 0 */2  * *   cd  /opt/BlackIP/blackwtredis && rm -rf blackip_list && git clone https://github.com/njcx/blackip_list.git  && bash /opt/BlackIP/blackwtredis/biptoredis.sh > /dev/null 2>&1
